# shairport-sync-pulseaudio/etc/s6-overlay/s6-rc.d/05-pulseaudio/run
#!/command/with-contenv sh
echo "Starting pulseaudio"
exec /usr/bin/pulseaudio --system --disallow-exit --disallow-module-loading
